# supervisor/worker.conf
[program:{{domain}}-worker]
process_name={{domain}}-worker
command=php {{deploy_path}}/current/artisan queue:work
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile={{deploy_path}}/log/worker.log
stopwaitsecs=216000
